# roles/letsencrypt/files/certbot.service
[Service]
ExecStart=/usr/bin/certbot --nginx -q renew